# services/reranker/baseline/deploy/reranker-vllm.service
[Unit]
Description=vLLM server for reranker (Qwen3.5-27B)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.reflection/services/reranker
ExecStart=%h/.reflection/services/reranker/.venv/bin/vllm serve Qwen/Qwen3.5-27B --port 42984 --dtype float16 --max-model-len 4096
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
