# services/reranker/baseline/deploy/reranker.service
[Unit]
Description=Reranker FastAPI server (cross-encoder via SGLang)
After=network.target reranker-baseline.service

[Service]
Type=simple
WorkingDirectory=%h/.reflection/services/reranker
ExecStart=%h/.reflection/services/reranker/.venv/bin/python -m services.reranker.baseline --host 0.0.0.0 --port 42983
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
